# core/darkmatter-worker.service
[Unit]
Description=DarkMatterK3@Home T4 Worker
After=network.target

[Service]
User=your_username
WorkingDirectory=/home/your_username/DarkMatterK3-Home/core
EnvironmentFile=/home/your_username/DarkMatterK3-Home/core/.env
ExecStart=/home/your_username/DarkMatterK3-Home/core/venv/bin/python /home/your_username/DarkMatterK3-Home/core/t4_worker.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
